# init/kcdnotif.debian.init.in
#!/bin/sh
### BEGIN INIT INFO
# Provides:          kcdnotif
# Required-Start:    $remote_fs $syslog
# Required-Stop:     $remote_fs $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Teambox Collaboration Daemon, Notification mode
# Description:       Teambox Collaboration Daemon, Notification mode
### END INIT INFO

PATH=/sbin:/usr/sbin:/bin:/usr/bin
DESC="Teambox collaboration daemon notification mode"
NAME=kcdnotif
DAEMON="@@PREFIX@@/bin/kcd"
SCRIPTNAME=/etc/init.d/$NAME.sh
LOCKFILE=/var/lock/kcdnotif.lock

case "$1" in
start)
        if [ `$KCD -P $LOCKFILE query` != "running" ]; then
            $KCD -P $LOCKFILE closefd start frontend
        fi
        if [ `$KCD -P $LOCKFILE query` != "running" ]; then
            log_failure_msg "Starting $DESC"
            exit 1
        else
            log_success_msg "Starting $DESC"
        fi
;;
status)
        if [ `$KCD -P $LOCKFILE query` == "running" ]; then
            echo "Running"
            exit 0
        else
            echo "Not running"
            exit 3
        fi
;;
stop)
        if [ `$KCD -P $LOCKFILE query` == "running" ]; then
            $KCD -P $LOCKFILE stop
        fi
        if [ `$KCD -P $LOCKFILE query` != "running" ]; then
            log_success_msg "Stopping $DESC"
            exit 0
        else
            log_failure_msg "Stopping $DESC"
            exit 1
        fi
;;

restart)
        $0 stop
  	$0 start
;;

*)
        echo "Usage: $0 {status|start|stop|restart}"
        exit 1
esac
